# bin/.config/systemd/user/check-mise-upgrades.service
[Unit]
Description=Check mise upgrades and send notification

[Service]
Type=oneshot
ExecStart=%h/bin/check-mise-upgrades
