# tmux/.config/systemd/user/weather.service
[Unit]
Description=Fetch weather stats

[Service]
Type=oneshot
ExecStart=%h/bin/weather.sh

[Install]
WantedBy=default.target
